# packaging/systemd/trollbridge.service
[Unit]
Description=trollbridge LLM-aware HTTP/HTTPS proxy
Documentation=https://github.com/dandriscoll/trollbridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=trollbridge
Group=trollbridge
ExecStart=/usr/local/bin/trollbridge run --config /etc/trollbridge/trollbridge.yaml
ExecReload=/bin/kill -HUP $MAINPID
# Operational log level: uncomment to override the default (info).
# Valid values: debug | info | warn | error.
# Equivalent to `--log-level` on the command line; read at startup
# (SIGHUP does not currently re-read it).
# Environment="TROLLBRIDGE_LOG_LEVEL=info"
KillSignal=SIGTERM
TimeoutStopSec=45
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
ReadWritePaths=/var/log/trollbridge /etc/trollbridge

# CAP_NET_BIND_SERVICE only needed when binding ports < 1024.
# AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
